# Munge-App/cron/crontab
MAILTO="matthijs@ischen.nl"
22,42 * * * *  ~/current/cron/run.pl ~/current/bin/sync.pl
