# Munge-App/cron/crontab
MAILTO="matthijs@ischen.nl"
0,30,45,15 * * * * PERL5LIB=/home/dotcloud/perl5/lib/perl5 /opt/perl5/perls/perl-5.12.4/bin/perl ~/current/bin/sync.pl
