# Munge-App/cron/crontab
MAILTO="matthijs@ischen.nl"
26,56 * * * * /opt/perl5/perls/current/bin/perl ~/current/bin/sync.pl
