# Munge-App/cron/crontab
MAILTO="matthijs@ischen.nl"
0,30,45,15 * * * *  ~/current/cron/run.pl ~/current/bin/sync.pl
